# roles/s3_log_sync/files/cron.daily
#!/bin/sh

/usr/local/bin/s3_sync.sh >> /var/log/s3_sync.log 2>&1
